# tools/whoop/whoop-refresh-monitor.service
[Unit]
Description=whoopskill token refresh + lightweight fetch

[Service]
Type=oneshot

Environment=PATH=/home/YOUR_USERNAME/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=WHOOP_CLIENT_ID=CLIENT_ID
Environment=WHOOP_CLIENT_SECRET=CLIENT_SECRET
Environment=WHOOP_REDIRECT_URI=https://youdomain.com/hooks/whoop

ExecStart=%h/.local/bin/whoop-refresh-monitor.sh
