# tools/whoop/whoop-refresh-monitor.timer
[Unit]
Description=Run whoopskill monitor periodically

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
